# Dockerfiles/source-destination-service-images/elasticsearch-kibana3/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/bin/bash -c "/usr/sbin/nginx"

[program:elasticsearch]
command=/bin/bash -c "/elasticsearch/bin/elasticsearch"
